# output/DEBIAN/postinst
#!/bin/bash
set -e

# Create ollama system user/group if they don't exist
if ! id -u ollama > /dev/null 2>&1; then
    adduser --system --shell /bin/false --group \
        --home /usr/share/ollama --no-create-home ollama
fi

# Reload systemd so the new service file is recognized
if command -v systemctl > /dev/null 2>&1 && systemctl is-system-running --quiet 2>/dev/null; then
    systemctl daemon-reload || true
fi

echo "Ollama installed. To complete setup:"
echo "  1. Start the service: sudo systemctl start ollama"
echo "  2. Run a model: ollama run llama3.2"
echo "  3. Enable on boot: sudo systemctl enable ollama"

# debian.griffo.io subscription notice (remove after the 1 Oct 2026 rollout settles)
if [ "$1" = "configure" ] && [ ! -e /etc/apt/auth.conf.d/debian.griffo.io.conf ]; then
    echo ""
    echo "debian.griffo.io: starting 1 October 2026 (noon UK), downloading packages"
    echo "from this repository requires a yearly subscription."
    echo "  How to subscribe: https://debian.griffo.io"
    echo "  Former Buy Me a Coffee supporters get a free lifetime subscription via Discord."
    echo ""
fi

exit 0
